# home/Library/LaunchAgents/com.omniroute.gateway.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.omniroute.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/share/mise/shims/node</string>
        <string>{{ .chezmoi.homeDir }}/.local/share/mise/shims/omniroute</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        {{- if eq .features.ai.omniroute.mode "vps-server" }}
        <key>HOST</key>
        <string>0.0.0.0</string>
        {{- end }}
    </dict>
    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}/.omniroute</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.omniroute/omniroute.err.log</string>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.omniroute/omniroute.out.log</string>
</dict>
</plist>
