# home/private_dot_config/private_systemd/private_user/hermes-dashboard.service.tmpl
# =====================================================================
# HERMES DASHBOARD — Systemd User Service
# =====================================================================
[Unit]
Description=Hermes Agent Dashboard
After=network.target

[Service]
Type=simple
Environment=HOME=%h
EnvironmentFile=%h/.hermes/.env
WorkingDirectory=%h
ExecStart=%h/.local/bin/hermes dashboard --port 9119 --host 0.0.0.0 --isolated
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=30
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target