# home/private_dot_config/private_systemd/private_user/omniroute.service.tmpl
[Unit]
Description=OmniRoute — LLM API Router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ .chezmoi.homeDir }}/.local/share/mise/installs/node/24/bin/omniroute
Restart=on-failure
RestartSec=5
Environment=PATH={{ .chezmoi.homeDir }}/.local/bin:{{ .chezmoi.homeDir }}/.local/share/mise/installs/node/24/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
