# home/private_dot_config/private_systemd/private_user/omniroute.service.tmpl
[Unit]
Description=OmniRoute — Smart AI Router with Auto Fallback
Documentation=https://github.com/omniroute/omniroute
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ .chezmoi.username }}
Group={{ .chezmoi.username }}
WorkingDirectory={{ .chezmoi.homeDir }}/.local/src/OmniRoute

# Node.js via mise (source build)
ExecStart={{ .chezmoi.homeDir }}/.local/share/mise/installs/node/lts/bin/node {{ .chezmoi.homeDir }}/.local/src/OmniRoute/bin/omniroute.mjs serve

# Environment
EnvironmentFile={{ .chezmoi.homeDir }}/.omniroute/.env
Environment=PATH={{ .chezmoi.homeDir }}/.local/share/mise/installs/node/lts/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production

# Restart behaviour — keep it alive
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Safety
TimeoutStopSec=30
KillMode=process

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
