# home/private_dot_config/systemd/user/camofox.service.tmpl
[Unit]
Description=Camofox — Anti-detection browser (Camoufox/Firefox)
Documentation=https://github.com/jo-inc/camofox-browser
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=NODE_ENV=production
Environment=CAMOFOX_PORT=3030
ExecStart=%h/.local/share/mise/shims/node %h/camofox/node_modules/camofox-browser/dist/src/server.js
Restart=on-failure
RestartSec=10
WorkingDirectory=%h/camofox

[Install]
WantedBy=default.target
