# home/private_dot_config/systemd/user/hermes-gateway.service.tmpl
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart={{ .chezmoi.homeDir }}/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run
WorkingDirectory={{ .chezmoi.homeDir }}/.hermes
Environment="PATH={{ .chezmoi.homeDir }}/.hermes/hermes-agent/venv/bin:{{ .chezmoi.homeDir }}/.hermes/hermes-agent/node_modules/.bin:/usr/bin:{{ .chezmoi.homeDir }}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV={{ .chezmoi.homeDir }}/.hermes/hermes-agent/venv"
Environment="HERMES_HOME={{ .chezmoi.homeDir }}/.hermes"
Restart=always
RestartSec=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
TimeoutStopSec=210
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
